# docker/php/horizon.conf
[program:horizon]
command=php /var/www/artisan horizon
autostart=true
autorestart=true
stderr_logfile=/var/www/storage/logs/horizon.err.log
stdout_logfile=/var/www/storage/logs/horizon.log
